# host/spectre-agent.service
[Unit]
Description=Spectre Host Agent
After=network.target

[Service]
Type=simple
User=root
# Adjust WorkingDirectory and ExecStart to the actual path where you clone or copy
# the panel repository on your host machine (Linux).
# For example, if you place files in /opt/spectre/host:
WorkingDirectory=/opt/spectre/host
ExecStart=/usr/bin/python3 /opt/spectre/host/host_agent.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
